# .claude/launchd/com.claude.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>{{HOME}}/.cloudflared/config.yml</string>
        <string>run</string>
        <string>{{TUNNEL_NAME}}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{PROJECT_DIR}}/.claude/logs/tunnel-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_DIR}}/.claude/logs/tunnel-stderr.log</string>
</dict>
</plist>
