# .claude/launchd/com.claude.webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.webhook</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jinguomini/server/pr-resolver/.claude/webhook/server.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WEBHOOK_CONFIG</key>
        <string>/Users/jinguomini/server/pr-resolver/.claude/webhook/config.toml</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/jinguomini/.local/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/jinguomini/server/pr-resolver/.claude/webhook</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jinguomini/server/pr-resolver/.claude/logs/webhook-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jinguomini/server/pr-resolver/.claude/logs/webhook-stderr.log</string>
</dict>
</plist>
